# scripts/supervisor_lvore.ini
[program:localvore]
command=/home/www/.local/bin/gunicorn localvore.wsgi
numprocs=1
directory=/home/www/Localvore/django-app
stdout_logfile=/home/www/Localvore/logs/server.log
redirect_stderr=true
autostart=true
autorestart=true
user=www
